# services/teambot.service
[Unit]
Description=Discord Scrim team bot

[Service]
Type=simple
Restart=always
Environment="TOKEN="
ExecStart=/usr/local/bin/teambot

[Install]
WantedBy=multi-user.target
